# rigctld/rigctld.service
[Unit]
Description=rigctld Hamradio rig controler
After=syslog.target network.target

[Service]
Type=simple
# rigctld flags:
#   -m 3070       Hamlib rig model (3070 = Icom IC-7100; list others with `rigctl -l`)
#   -r /dev/...   serial device for CI-V/CAT control
#   -t 4532       TCP port rigctld listens on (lockstep, WSJT-X, etc. connect here)
#   -s 19200       serial speed; MUST match the radio's CI-V baud setting
# Tip: /dev/ttyUSB0 can change with USB enumeration order. For a stable path use
#      a /dev/serial/by-id/... symlink (see: ls -l /dev/serial/by-id/).
ExecStart=/usr/bin/rigctld -m 3070 -r /dev/ttyUSB0 -t 4532 -s 19200
ExecReload=/bin/kill -HUP $MAINPID
RestartSec=60
Restart=always
User=rigctld
Group=rigctld

[Install]
WantedBy=multi-user.target
